# roles/s3-ssh-keys/files/update-keys.service
[Unit]
Description=Update SSH keys

[Service]
ExecStart=/opt/features/ssh-keys/install-keys.sh

[Install]
WantedBy=multi-user.target
