# roles/airflow/templates/profile.d/airflow_python_libs.sh.j2
#!/usr/bin/env bash

export PYTHONPATH="$PYTHONPATH:{{ airflow_python_libs }}"
