# roles/github-ssh-keys/files/update-keys.timer
[Unit]
Description=Update SSH keys every half hour

[Timer]
OnBootSec=30 seconds
OnUnitActiveSec=30 minutes

[Install]
WantedBy=timers.target